# mxversions.service
[Unit]
Description=mxversions
[Service]
ExecStart=/root/mxversions/reportloop.sh

[Install]
WantedBy=multi-user.target
